# ansible/templates/payreconcile.service.j2
[Unit]
Description=PayReconcile Backend Service
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/opt/payreconcile/backend
Environment=PORT=5000
Environment=STRIPE_SECRET_KEY={{ stripe_secret }}
Environment=STRIPE_PUBLISHABLE_KEY={{ stripe_publishable }}
ExecStart=/usr/bin/java -jar app.jar
Restart=always
RestartSec=10
StandardOutput=append:/var/log/payreconcile.log
StandardError=append:/var/log/payreconcile.log

[Install]
WantedBy=multi-user.target
